# root/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-change-shell/up
foreground { echo "Running init-mod-openssh-server-change-shell/run" }
/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-change-shell/run
