# root/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-tweaks/run
#!/usr/bin/with-contenv bash
echo "I am $0"
exit 0
