# root/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-tweaks/up
foreground { echo "Running init-mod-openssh-server-tweaks/run" }
/etc/s6-overlay/s6-rc.d/init-mod-openssh-server-tweaks/run
